# system/bt-reconnect.timer
[Unit]
Description=Run Bluetooth reconnect every 30 seconds

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
